# systemd/trading-bot.service
[Unit]
Description=Trading Bot - Telegram Signal Processor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/trading_bot
Environment="PATH=/opt/trading_bot/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Main process
ExecStart=/opt/trading_bot/.venv/bin/python main.py

# Restart only on failure (crashes), not on normal stop
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Resource limits
MemoryMax=1G
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-bot

[Install]
WantedBy=multi-user.target
