# systemd/trading-bot-stop.timer
[Unit]
Description=Stop Trading Bot at 3:35 PM (Mon-Fri)
Requires=trading-bot.service

[Timer]
# Run at 3:35 PM on weekdays (Mon-Fri)
OnCalendar=Mon-Fri 15:35:00
Persistent=true

[Install]
WantedBy=timers.target
