# systemd/trading-bot-start.timer
[Unit]
Description=Start Trading Bot at 9:00 AM (Mon-Fri)
Requires=trading-bot.service

[Timer]
# Run at 9:00 AM on weekdays (Mon-Fri)
OnCalendar=Mon-Fri 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
